# config/movie.conf
[program:movie]
command=/home/john/venv/bin/gunicorn django_movie.wsgi:application -c /home/john/django_movie/config/gunicorn.conf.py
directory=/home/john/django_movie
user=john
autorestart=true
redirect_stderr=true
stdout_logfile = /home/john/django_movie/logs/debug.log
